# helloworld-api.service
[Unit]
Description=helloworld-api

[Service]
Environment="HELLOWORLD_API_ENV=PROD"
WorkingDirectory=/{{HOME}}/helloworld-api/
ExecStart=/{{HOME}}/go/bin/helloworld-api
User=root
Group=root
Restart=always

[Install]
WantedBy=multi-user.target
